# files/opt/pods/init.d/pod-proxy
#!/sbin/openrc-run

command_user="alex"

pod_name="proxy"
pod_networks="proxy media"
pod_yaml="/home/${command_user}/pods/${pod_name}.yml"

source "/home/${command_user}/pods/init.d/pod"
