# voidrepo/crontab
# redirects stdout and stderr to docker logs, view with `docker compose logs -f`
0 */3 * * * /usr/local/bin/sync.sh >> /proc/1/fd/1 2>> /proc/1/fd/2

